# deploy/systemd/flow-task-web.service
[Unit]
Description=Flow Task Frontend Service
After=network.target

[Service]
Type=simple
User=dudewei
Group=dudewei
WorkingDirectory=/home/dudewei/projects/comfyUI-flow-task/frontend
Environment=NODE_ENV=production
ExecStartPre=/usr/bin/env npm run build
ExecStart=/usr/bin/env npm run preview -- --host 0.0.0.0 --port 5173
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
